# backend/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F

[program:unit]
command=/usr/sbin/unitd